# recipes-test/simpleechoserver/simpleechoserver/simpleechoserver.service
[Unit]
Description=A simple echo server

[Service]
ExecStart=@BINDIR@/simple_echo_server 12345
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target